# bot.conf
[program:bot]
command=node-6 /home/leo/bot/src/bot.js
autostart=true
autorestart=true
stderr_logfile=/home/leo/bot/log/err.log
stdout_logfile=/home/leo/bot/log/out.log
